# flask-app.service
[Unit]
Descriotion=Flask app service
After=network.target

[Service]
WorkingDirectory=/home/maltish/app
ExecStart=/home/maltish/app/venv/bin/python main.py

[Install]
WantedBy=multi-user.target
